# systemd/trick-skel.timer
[Unit]
Description=Run trick_skel job daily

[Timer]
OnCalendar=daily
